# root/etc/services.d/activemq/finish
#!/usr/bin/execlineb

s6-svscanctl -t /var/run/s6/services
